# config/user.chronometry.webserver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>user.chronometry.webserver</string>
    
    <!-- Program to run -->
    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_ROOT}}/venv/bin/python</string>
        <string>{{PROJECT_ROOT}}/src/web_server.py</string>
    </array>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>{{PROJECT_ROOT}}</string>
    
    <!-- Run at load (boot) -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep alive - restart if crashed -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <!-- Restart after 10 seconds if crashed -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>{{PROJECT_ROOT}}/logs/webserver.log</string>
    
    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>{{PROJECT_ROOT}}/logs/webserver.error.log</string>
    
    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
</dict>
</plist>

